# overmind.service
[Unit]
Description=overmind for bug hexapod

[Service]
ExecStart=/root/stable/overmind

[Install]
WantedBy=multi-user.target
