# scripts/overmind.service
[Unit]
Description=overmind for bug hexapod

[Service]
ExecStart=/home/root/start.sh

[Install]
WantedBy=multi-user.target
